# x11/autostart/xkb-layout-monitor.desktop
[Desktop Entry]
Exec=sh -c "\$HOME/.config/dbus/xkb-layout-monitor.py"
Name[en_US]=XKB Layout Monitor
Type=Application
