# x11/autostart/xmodmap.desktop
[Desktop Entry]
Exec=sh -c "xmodmap \$HOME/.Xmodmap &"
Name[en_US]=Xmodmap
Type=Application
